# autonomy-perception-2d-main/services/autonomy-perception_2d.service
[Unit]
RequiresMountsFor=/home
After=network.target multi-user.target
Description=2d_perception_service

[Service]
type=notify
User=linkxavier
Environment="RMW_IMPLEMENTATION=rmw_cyclonedds_cpp"
Environment="CYCLONEDDS_URI=file:///home/linkxavier/robot_cfg/cyclonedds/cyclonedds.xml"
ExecStart=/home/linkxavier/scripts/autonomy-perception-2d.sh
Restart=on-failure
WatchdogSec=0
NotifyAccess=main

[Install]
WantedBy=multi-user.target
